# zeta/daphneservice.service
[Service]
ExecStart=/root/daphneboot.sh

[Install]
WantedBy=default.target
